# eunos_overlay/etc/local.d/network.start
#!/bin/sh
# Auto network setup
for iface in $(ip link show 2>/dev/null | grep -E "^[0-9]+:" | awk -F': ' '{print $2}' | grep -v lo | grep -v wlan); do
    ip link set "$iface" up 2>/dev/null
    ip addr add 10.0.2.15/24 dev "$iface" 2>/dev/null
    ip route add default via 10.0.2.2 2>/dev/null
    echo "nameserver 8.8.8.8" > /etc/resolv.conf
    break
done
